# board/rpi-client/rootfs_overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -o remount,rw /
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/bin/mount -r /dev/mmcblk0p1 /boot
::sysinit:/bin/mount -a
::sysinit:/bin/snapcastbox.sh
# ::sysinit:/bin/hostname -F /etc/hostname
# Hostname set by snapcastbox.sh

# now run any rc scripts
::sysinit:/etc/init.d/rcS

# Start Remote Control Programm
#tty1::respawn:/bin/snapremote.py

# Put a getty on the serial port
#ttyAMA0::respawn:/sbin/getty -L  console 0 vt100 # GENERIC_SERIAL
tty1::respawn:/sbin/getty -L  tty1 0 vt100 # Login console on tty1
#tty2::respawn:/sbin/getty -L  tty2 0 vt100 # Login console on tty2

# Stuff to do for the 3-finger salute
#::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
